# srcs/ftps/supervisord.conf
[supervisord]
nodaemon=true

[program:vsftpd]
command = sh -c "vsftpd /etc/vsftpd/vsftpd.conf || kill -s SIGTERM $(cat supervisord.pid)"
exitcode = 137

[program:telegraf]
command = sh -c "telegraf && kill -s SIGTERM $(cat supervisord.pid)"
